# site-template/systemd/epicsarchiverap.service.in
[Unit]
Description=EPICS Archiver Appliance for @ARCHAPPL_SITEID@
Documentation=@DOCURL@
After=network.target
SourcePath=@INSTALL_LOCATION@/@ARCHAPPL_MAIN_SCRIPT@

[Service]
#
User=@USERID@
Group=@GROUPID@
#
ExecStart=/bin/bash -c "@INSTALL_LOCATION@/@ARCHAPPL_MAIN_SCRIPT@ startup"
ExecStop=/bin/bash -c "@INSTALL_LOCATION@/@ARCHAPPL_MAIN_SCRIPT@ shutdown"
Type=forking

[Install]
WantedBy=multi-user.target
Alias=archappl.service
