# site-template/systemd/tomcat@mgmt.service.in
[Unit]
After=mariadb.service
Requires=mariadb.service
#PartOf=@APPNAME@.service

[Service]
User=@USERID@
Group=@GROUPID@

Environment="ARCHAPPL_APPLIANCES=@INSTALL_LOCATION@/@ARCHAPPL_APPLIANCES@"

Environment="JAVA_HOME=@JAVA_HOME@"
Environment="JAVA_OPTS=@JAVA_OPTS@"

Environment="CATALINA_BASE=@TOMCAT_HOME@"
Environment="CATALINA_HOME=@INSTALL_LOCATION@/@SERVICE@"
Environment="CATALINA_OPTS=@CATALINA_OPTS@"

Type=forking

ExecStart=@INSTALL_LOCATION@/@SERVICE@/bin/startup.sh
ExecStop=@INSTALL_LOCATION@/@SERVICE@/bin/shutdown.sh

[Install]
WantedBy=multi-user.target

