# site-template/cf.service.in
[Unit]
Description=ChannelFinder Service
Documentation=https://github.com/ChannelFinder/ChannelFinder-SpringBoot
After=network.target elasticsearch.service
Requires=elasticsearch.service

[Service]
# User=elasticsearch
# Group=elasticsearch

ExecStart=_JAVAPATH_/java _JAVAOPTS_ -jar _CFPATH_/_CHANNELFINDER_JAR_NAME_
SuccessExitStatus=143 

ExecReload=kill -SIGINT $MAINPID
KillMode=process

Restart=on-failure
RestartSec=10s

[Install]
WantedBy=multi-user.target
Alias=cf.service