# site-template/systemd.service.in
# After changes
#   sudo systemctl daemon-reload
#
# Start/status/stop
#   sudo systemctl start  @TARGET_NAME@.service
#   sudo systemctl status @TARGET_NAME@.service
#   sudo systemctl stop   @TARGET_NAME@.service
#
# Start when computer boots:
#   sudo systemctl enable @TARGET_NAME@.service
#

[Unit]
Description=@TARGET_NAME@ Service
Documentation=@DOCURL@
After=syslog.target network.target remote-fs.target

[Service]
Environment=AS_IMQ_LIB="@INSTALL_LOCATION@/mq/lib"
Environment=AS_IMQ_BIN="@INSTALL_LOCATION@/mq/bin"
Environment=AS_CONFIG="@INSTALL_LOCATION@/config"
Environment=AS_INSTALL="@INSTALL_LOCATION@"
Environment=AS_DEF_DOMAINS_PATH="@INSTALL_LOCATION@/domains"
Environment=AS_DEF_NODES_PATH="@INSTALL_LOCATION@/nodes"
Environment=AS_DERBY_INSTALL="@INSTALL_LOCATION@/javadb"
Environment=AS_H2_INSTALL="@INSTALL_LOCATION@/h2db"

#
User=@USERID@
Group=@GROUPID@
#
Type=oneshot
RemainAfterExit=yes
#
ExecStart=@JAVA_PATH@/java -XX:+IgnoreUnrecognizedVMOptions -jar @INSTALL_LOCATION@/glassfish/lib/client/appserver-cli.jar start-domain production
ExecReload=@JAVA_PATH@/java -XX:+IgnoreUnrecognizedVMOptions -jar @INSTALL_LOCATION@/glassfish/lib/client/appserver-cli.jar restart-domain production
ExecStop=@JAVA_PATH@/java -XX:+IgnoreUnrecognizedVMOptions -jar @INSTALL_LOCATION@/glassfish/lib/client/appserver-cli.jar stop-domain production
#
TimeoutStartSec=300
TimeoutStopSec=30
#
[Install]
WantedBy = multi-user.target
#Alias=glassfish.service